# hexo.ini
[program:hexo]
command=/bin/hexo s
directory=/root/blog/blog.uzi.cool
