# hexo.ini
[program:hexo]
command=/bin/hexo s
directory=/app/web/blog.uzi.cool
